# etc/service/hydrator/run
#!/bin/sh
exec 2>&1
source /etc/envvars
echo " ---- starting hydrator "
/usr/bin/hydrator
if [ $? -ne 0 ]; then
  echo " ---- unable to start hydrator (sleeping 5 seconds) "
  sleep 5 
fi
